# bucketsync.timer
[Unit]
Description=efs %p of efs backup every 2 hours

[Timer]
OnBootSec=10m
OnUnitActiveSec=2h
Unit=bucketsync.service

